# com.toabh-casting-hub.flask-backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.toabh-casting-hub.flask-backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>flask</string>
        <string>run</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>5000</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/aineshsikdar/github-dashboards/toabh-casting-hub</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>FLASK_APP</key>
        <string>api/index.py</string>
        <key>PATH</key>
        <string>/Users/aineshsikdar/.hermes/hermes-agent/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/aineshsikdar/github-dashboards/toabh-casting-hub/backend/flask.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/aineshsikdar/github-dashboards/toabh-casting-hub/backend/flask.error.log</string>
</dict>
</plist>
